# grub-cleanup.service
[Unit]
Description=One-time GRUB cleanup after Clonezilla operation
ConditionPathExists=/etc/grub.d/40_custom

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/grub-cleanup.sh

[Install]
WantedBy=multi-user.target